# cloudbreak.conf
description "A job that starts cloudbreak"
author "Malathi"

start on runlevel [2345]
setuid ubuntu

script

    cd /home/ubuntu/cloudbreak
    cbd start

end script

